# script/power_check.service
[Unit]
Description=My power_check servicee
After=multi-user.target

[Service]
Type=idle
ExecStart=/usr/bin/python3 -u /home/pi/Power_check/power_check.py 2>&1 |  /home/pi/Power_check/power_check.log 

[Install]
WantedBy=multi-user.target
